# .debian/prerm.sh
#!/usr/bin/env bash

systemctl disable burstcoin-wallet || true
systemctl stop burstcoin-wallet || true
